# deploy/primer-empleado-ia-backup.service
[Unit]
Description=Backup CRM - Encuentra Tu Primer Empleado IA
Wants=primer-empleado-ia.service

[Service]
Type=oneshot
WorkingDirectory=/opt/primer-empleado-ia
ExecStart=/usr/bin/python3 /opt/primer-empleado-ia/backup_crm.py
User=primeria
Group=primeria
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/opt/primer-empleado-ia
