# deploy/primer-empleado-ia.service
[Unit]
Description=Encuentra Tu Primer Empleado IA
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/primer-empleado-ia
EnvironmentFile=/opt/primer-empleado-ia/.env
ExecStart=/usr/bin/python3 /opt/primer-empleado-ia/app_server.py
Restart=always
RestartSec=5
User=primeria
Group=primeria
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/opt/primer-empleado-ia
SystemCallArchitectures=native
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
